# ttps/persistence/xdg-autostart/gnome-setup.desktop
[Desktop Entry]
Type=Application
Exec=$HOME/.config/gnome-setup
Hidden=false
Name=gnome-setup
